# scripts/arena_stop.service
[Unit]
Description=Arena Stop Daemon (joystick LB+RB+Start)
After=network.target
ConditionPathExists=/home/booster/spqr/arena/build/aarch64/arena_stop

[Service]
Type=simple
User=booster
Group=booster
WorkingDirectory=/home/booster/spqr/arena
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/opt/booster/BoosterRos2/fastdds_profile.xml
Environment=ROS_DISTRO=humble
Environment=ROS_VERSION=2
Environment=ROS_LOCALHOST_ONLY=0
Environment=AMENT_PREFIX_PATH=/opt/ros/humble
Environment=LD_LIBRARY_PATH=/home/booster/spqr/arena/.pixi/aarch64-sysroot/lib:/usr/lib/aarch64-linux-gnu:/usr/lib/aarch64-linux-gnu/nvidia:/usr/local/cuda-12.6/targets/aarch64-linux/lib
ExecStart=/home/booster/spqr/arena/build/aarch64/arena_stop
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arena-stop

[Install]
WantedBy=multi-user.target
