# scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:booster-motion]
directory=/opt/booster
command=/opt/booster/booster-motion -mode sim -config configs/config_isaac.lua
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/booster-motion.log
stderr_logfile=/var/log/booster-motion.err

[program:simbridge]
command=/opt/ros/humble/bin/simbridge
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/simbridge.log
stderr_logfile=/var/log/simbridge.err

[program:arena]
directory=/workspace
command=/bin/bash --norc --noprofile -c "sleep 2 && exec /workspace/.pixi/envs/default/x86_64-conda-linux-gnu/sysroot/lib64/ld-linux-x86-64.so.2 --library-path /workspace/.pixi/envs/default/x86_64-conda-linux-gnu/sysroot/lib64:/opt/booster/lib-usr-local:/opt/booster/lib:/workspace/.pixi/envs/default/lib:/opt/booster/lib-x86_64-linux-gnu /workspace/build/x86/main --backend booster --task t1-velocity-flat"
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/arena.log
stderr_logfile=/var/log/arena.err
